# scripts/com.agentcos.sidecar.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.agentcos.sidecar</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/tonyseneadza/Codehome/AgenticOS/.venv/bin/python</string>
    <string>-m</string>
    <string>gui.sidecar</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/tonyseneadza/Codehome/AgenticOS</string>

  <key>StandardOutPath</key>
  <string>/Users/tonyseneadza/Codehome/AgenticOS/data/logs/sidecar.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/tonyseneadza/Codehome/AgenticOS/data/logs/sidecar.log</string>

  <!-- Restart automatically on crash; wait 10 s between attempts -->
  <key>KeepAlive</key>
  <true/>
  <key>ThrottleInterval</key>
  <integer>10</integer>

  <!-- Start when this agent is loaded (login / launchctl bootstrap) -->
  <key>RunAtLoad</key>
  <true/>

  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>/Users/tonyseneadza</string>
    <key>PATH</key>
    <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>
</dict>
</plist>
